# pick-dashboard-restart.service
[Unit]
Description=Restart Pick Dashboard service

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart pick-dashboard.service
